# docker/update.cron
0 * * * * PATH=$PATH:/usr/bin:/usr/local/bin /usr/local/bin/entrypoint update > /proc/1/fd/1 2>/proc/1/fd/2
